# rootfs/etc/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:tftpd]
command=exec in.tftpd -L -vvv -u ftp --secure --address 0.0.0.0:1069 "${TFTPD_EXTRA_ARGS}" /tftpboot